# custom-updater.timer
[Unit]
Description=Custom Updater Timer

[Timer]
OnUnitActiveSec=120s
OnBootSec=30s

[Install]
WantedBy=timers.target
